# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:django]
command=python manage.py runserver 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/django/django.log
stderr_logfile=/var/log/django/django-error.log
priority=1

[program:celery-worker]
command=celery -A oraclemanager worker --loglevel=info --concurrency=2
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/django/celery-worker.log
stderr_logfile=/var/log/django/celery-worker-error.log
priority=2

[program:celery-beat]
command=celery -A oraclemanager beat --loglevel=info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/django/celery-beat.log
stderr_logfile=/var/log/django/celery-beat-error.log
priority=3

[program:streamer]
command=python manage.py run_streamer
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/django/streamer.log
stderr_logfile=/var/log/django/streamer-error.log
priority=4

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
